# .config/systemd/user/calsync.timer
[Unit]
Description=Runs calsync perodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Unit=mailsync.service

[Install]
WantedBy=timers.target
