# .config/systemd/user/mailcalsync.service
[Unit]
Description=Runs outlooksync perodically

[Service]
Type=simple
ExecStart=/home/barts/bin/outlooksync

[Install]
WantedBy=default.target
