# .config/systemd/user/vdirsyncer.timer
[Unit]
Description=Synchronize vdirs

[Timer]
OnBootSec=1m
OnUnitActiveSec=15m
AccuracySec=5m

[Install]
WantedBy=timers.target

