# dotfiles-linux/systemd/user/hosts.timer
[Unit]
Description=Update hosts

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
