# dotfiles-linux/systemd/user/mailsync.timer
[Unit]
Description=Runs mailsync perodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Unit=mailsync.service

[Install]
WantedBy=timers.target
