# dotfiles-linux/systemd/user/system-check.service
[Unit]
Description=Run Arch system health check

[Service]
Type=oneshot
Environment=SYSTEM_CHECK_SKIP_SUDO=1
ExecStart=/home/barts/dotfiles-linux/codex/skills/system-check/scripts/run-scheduled-system-check.sh
