# dotfiles-linux/systemd/root/ipset-load.service
[Unit]
Description=Restore ipset sets at boot
Before=iptables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ipset-load blocked /etc/ipset.conf

[Install]
WantedBy=multi-user.target

