# dotfiles-linux/systemd/root/update-blocklist.service
[Unit]
Description=Update IP blocklists

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-blocklist

